# packaging/systemd/praetor-web.service
[Unit]
Description=Praetor shared web client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/praetor/web.env
ExecStart=%h/.local/bin/praetor-web --listen 127.0.0.1:8787
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=0077

[Install]
WantedBy=default.target
